# laravel-hhvm/container/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * hhvm /data/www/artisan schedule:run >> /dev/null 2>&1
