# laravel-hhvm/container/crontab
* * * * * hhvm /data/www/artisan schedule:run
